# crontab
0 9 * * * /usr/local/bin/python /usr/src/app/main.py >> /var/log/cron.log 2>&1
# Don't forget to end with an empty newline!
